# profile.d/01_alias.sh
alias ..="cd .."
alias ...="cd ../.."
alias lh='ls -lhStA --time-style="+${DATE_TIME_FORMAT}"'
alias now='date "+${DATE_TIME_FORMAT}"'
alias watch="watch -d -n1"
alias du1="du -ach -d1 -x|sort -h"
alias du2="du -ach -d2 -x|sort -h"
alias dfh="df -Th"
alias freeh="free -mwh"
alias psl="ps auxf"
alias psearch="ps aux | grep -v grep | grep -i -e VSZ -e"
